# famcal.service
[Unit]
Description=Family Calendar ICS Server
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/famcal
Environment="PATH=/home/pi/famcal/.venv/bin"
Environment="FAMCAL_USE_DATABASE=true"
Environment="FAMCAL_DATABASE_ICS=true"
Environment="FAMCAL_API_MEMBER_EVENTS_LIMIT=100000"
Environment="FAMCAL_API_COMBINED_EVENTS_LIMIT=1000000"
Environment="SQLALCHEMY_DATABASE_URL=postgresql://famcal_user:YOUR_PASSWORD@localhost:5432/famcal_db"
Environment="FAMCAL_WEB_PASSWORD=YOUR_SITE_PASSWORD"
Environment="FAMCAL_SECRET_KEY=YOUR_RANDOM_SECRET_KEY"
ExecStart=/home/pi/famcal/.venv/bin/python /home/pi/famcal/family_calendar_server.py --config /home/pi/famcal/family_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
